# systemd/mp3pi.service
[Unit]
Description=mp3pi Kivy App

[Service]
#Type=forking
ExecStart=/root/mp3pi/startup.sh

[Install]
WantedBy=multi-user.target
